# zabbix-quanta/zabbix-quanta-iostats-1.0.0/debian/postinst
#!/bin/sh -e

if which invoke-rc.d >/dev/null 2>&1; then
  invoke-rc.d zabbix-agent restart
else
  /etc/init.d/zabbix-agent restart
fi

#DEBHELPER#

exit 0
